# supervisor/0-install.conf
[program:install]
command=/bin/bash -c "/app/scripts/startup.sh"
numprocs=1
startsecs=0
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0