# supervisor/1-cron.conf
[program:cron]
command=/usr/sbin/cron -f
autostart = true
autorestart=true
priority=5