# supervisor/3-php.conf
[program:php]
command = /usr/sbin/php-fpm7.4 -R -F
autostart = true
autorestart = true
priority=2
