# supervisor/3-php.conf
[program:php]
command = /usr/sbin/php-fpm8.1 -R -F
autostart = true
autorestart = true
priority=2
