# supervisor/3-php.conf
[program:php]
command = /etc/init.d/php7.2-fpm start
autostart = true
autorestart=true
priority=2