# supervisor/3-php.conf
[program:php]
command = /usr/sbin/php5-fpm
autostart = true
autorestart=true
priority=2