# supervisor/4-node.conf
[program:node]
command=/usr/bin/node /app/server.js
autostart=true
autorestart=true
priority=4
